# pkg-installer/scripts/postinstall
#!/bin/bash
# BriefDesk Post-Install Script
# This runs as root after pkg installation, but we need to set up user-space components

set -e

# Get the user who invoked the installer (not root)
INSTALL_USER="${USER:-$(stat -f '%Su' /dev/console)}"
USER_HOME=$(eval echo ~$INSTALL_USER)

# Source files installed by pkg to /usr/local/share/briefdesk
PKG_SOURCE="/usr/local/share/briefdesk"

# User directories
INSTALL_DIR="$USER_HOME/.local/share/briefdesk"
DEVSAI_DIR="$USER_HOME/.local/share/devsai"
LAUNCHAGENTS_DIR="$USER_HOME/Library/LaunchAgents"

echo "Setting up BriefDesk for user: $INSTALL_USER"

# Create directories
sudo -u "$INSTALL_USER" mkdir -p "$INSTALL_DIR"
sudo -u "$INSTALL_USER" mkdir -p "$DEVSAI_DIR"
sudo -u "$INSTALL_USER" mkdir -p "$LAUNCHAGENTS_DIR"

# Copy application files
cp -R "$PKG_SOURCE/"* "$INSTALL_DIR/"
chown -R "$INSTALL_USER" "$INSTALL_DIR"

# Find Python (with extensive fallbacks for different installation methods)
SYSTEM_PYTHON=""

# Helper function to check if Python works and is recent enough (3.9+)
check_python() {
    local py="$1"
    if [ -x "$py" ] && "$py" -c "import sys; exit(0 if sys.version_info >= (3, 9) else 1)" 2>/dev/null; then
        return 0
    fi
    return 1
}

# 1. Check Homebrew on Apple Silicon
if [ -z "$SYSTEM_PYTHON" ] && [ -f "/opt/homebrew/bin/python3" ]; then
    check_python "/opt/homebrew/bin/python3" && SYSTEM_PYTHON="/opt/homebrew/bin/python3"
fi

# 2. Check versioned Homebrew Python on Apple Silicon (prefer newer versions)
if [ -z "$SYSTEM_PYTHON" ]; then
    for ver in 3.14 3.13 3.12 3.11 3.10 3.9; do
        if [ -f "/opt/homebrew/bin/python$ver" ]; then
            check_python "/opt/homebrew/bin/python$ver" && SYSTEM_PYTHON="/opt/homebrew/bin/python$ver" && break
        fi
    done
fi

# 3. Check Python.org installer (Framework install)
if [ -z "$SYSTEM_PYTHON" ]; then
    for ver in 3.14 3.13 3.12 3.11 3.10 3.9; do
        if [ -f "/Library/Frameworks/Python.framework/Versions/$ver/bin/python3" ]; then
            check_python "/Library/Frameworks/Python.framework/Versions/$ver/bin/python3" && \
                SYSTEM_PYTHON="/Library/Frameworks/Python.framework/Versions/$ver/bin/python3" && break
        fi
    done
fi

# 4. Check Homebrew on Intel Macs
if [ -z "$SYSTEM_PYTHON" ] && [ -f "/usr/local/bin/python3" ]; then
    check_python "/usr/local/bin/python3" && SYSTEM_PYTHON="/usr/local/bin/python3"
fi

# 5. Check versioned Homebrew Python on Intel Macs
if [ -z "$SYSTEM_PYTHON" ]; then
    for ver in 3.14 3.13 3.12 3.11 3.10 3.9; do
        if [ -f "/usr/local/bin/python$ver" ]; then
            check_python "/usr/local/bin/python$ver" && SYSTEM_PYTHON="/usr/local/bin/python$ver" && break
        fi
    done
fi

# 6. Check pyenv (common for developers)
if [ -z "$SYSTEM_PYTHON" ] && [ -f "$HOME/.pyenv/shims/python3" ]; then
    check_python "$HOME/.pyenv/shims/python3" && SYSTEM_PYTHON="$HOME/.pyenv/shims/python3"
fi

# 7. Check system Python (macOS Sonoma+ has Python 3.9+)
if [ -z "$SYSTEM_PYTHON" ] && [ -f "/usr/bin/python3" ]; then
    check_python "/usr/bin/python3" && SYSTEM_PYTHON="/usr/bin/python3"
fi

# 8. Last resort: use which
if [ -z "$SYSTEM_PYTHON" ]; then
    FOUND_PYTHON=$(which python3 2>/dev/null)
    if [ -n "$FOUND_PYTHON" ] && check_python "$FOUND_PYTHON"; then
        SYSTEM_PYTHON="$FOUND_PYTHON"
    fi
fi

# Final fallback (may not work but we try)
if [ -z "$SYSTEM_PYTHON" ]; then
    echo "WARNING: Could not find a suitable Python 3.9+ installation"
    SYSTEM_PYTHON="/usr/bin/python3"
fi

echo "Using Python: $SYSTEM_PYTHON ($($SYSTEM_PYTHON --version 2>&1))"

# Create Python wrapper
cat > "$INSTALL_DIR/python3" << PYEOF
#!/bin/bash
exec "$SYSTEM_PYTHON" "\$@"
PYEOF
chmod +x "$INSTALL_DIR/python3"
chown "$INSTALL_USER" "$INSTALL_DIR/python3"

# Ensure pip is available (bootstrap if needed)
echo "Checking pip availability..."
if ! sudo -u "$INSTALL_USER" "$SYSTEM_PYTHON" -m pip --version >/dev/null 2>&1; then
    echo "pip not found, bootstrapping..."
    sudo -u "$INSTALL_USER" "$SYSTEM_PYTHON" -m ensurepip --user 2>/dev/null || true
fi

# Install Google API libraries for OAuth
echo "Installing Google API libraries..."
# Try multiple pip install methods (some systems need --break-system-packages, some don't)
if sudo -u "$INSTALL_USER" "$SYSTEM_PYTHON" -m pip install --user google-auth-oauthlib google-api-python-client 2>&1; then
    echo "Google API libraries installed successfully"
elif sudo -u "$INSTALL_USER" "$SYSTEM_PYTHON" -m pip install --break-system-packages --user google-auth-oauthlib google-api-python-client 2>&1; then
    echo "Google API libraries installed successfully (with --break-system-packages)"
else
    echo "WARNING: Could not install Google API libraries. You may need to run manually:"
    echo "  pip3 install --user google-auth-oauthlib google-api-python-client"
fi

# Find Node.js
NODE_PATH=""
if [ -f "/usr/local/bin/node" ]; then
    NODE_PATH="/usr/local/bin/node"
elif [ -f "/opt/homebrew/bin/node" ]; then
    NODE_PATH="/opt/homebrew/bin/node"
elif [ -f "/opt/homebrew/opt/node/bin/node" ]; then
    NODE_PATH="/opt/homebrew/opt/node/bin/node"
fi

# Find npm (for devsai install)
NPM_PATH=""
if [ -f "/usr/local/bin/npm" ]; then
    NPM_PATH="/usr/local/bin/npm"
elif [ -f "/opt/homebrew/bin/npm" ]; then
    NPM_PATH="/opt/homebrew/bin/npm"
elif [ -n "$NODE_PATH" ] && [ -f "$(dirname "$NODE_PATH")/npm" ]; then
    NPM_PATH="$(dirname "$NODE_PATH")/npm"
else
    NPM_PATH="$(command -v npm 2>/dev/null || true)"
fi

# Set up Node for devsai if found
if [ -n "$NODE_PATH" ]; then
    cp "$NODE_PATH" "$DEVSAI_DIR/node"
    chmod +x "$DEVSAI_DIR/node"
    chown "$INSTALL_USER" "$DEVSAI_DIR/node"
fi

# Install devsai CLI from bundled tarball (if available)
DEVSAI_TARBALL="$INSTALL_DIR/devsai/devsai.tgz"
if [ -f "$DEVSAI_TARBALL" ] && [ -n "$NPM_PATH" ]; then
    echo "Installing devsai CLI from bundled tarball..."
    sudo -u "$INSTALL_USER" "$NPM_PATH" install --prefix "$DEVSAI_DIR" "$DEVSAI_TARBALL" --silent 2>/dev/null || \
        echo "WARNING: devsai install failed. You may need to reinstall devsai manually."

    if [ ! -f "$DEVSAI_DIR/node_modules/devsai/dist/index.js" ]; then
        echo "WARNING: devsai install incomplete (missing dist/index.js). Verify npm install and rerun."
    fi

    # Create devsai wrapper (user-local)
    cat > "$DEVSAI_DIR/devsai.sh" << 'EOF'
#!/bin/bash
exec "$HOME/.local/share/devsai/node" "$HOME/.local/share/devsai/node_modules/devsai/dist/index.js" "$@"
EOF
    chmod +x "$DEVSAI_DIR/devsai.sh"
    chown "$INSTALL_USER" "$DEVSAI_DIR/devsai.sh"

    # Add a convenient executable in ~/.local/bin
    sudo -u "$INSTALL_USER" mkdir -p "$USER_HOME/.local/bin"
    sudo -u "$INSTALL_USER" ln -sf "$DEVSAI_DIR/devsai.sh" "$USER_HOME/.local/bin/devsai"
else
    echo "WARNING: devsai tarball or npm not found. devsai CLI will not be installed."
fi

# Create LaunchAgent for static server
cat > "$LAUNCHAGENTS_DIR/com.briefdesk.static.plist" << EOF
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.briefdesk.static</string>
    <key>ProgramArguments</key>
    <array>
        <string>$INSTALL_DIR/python3</string>
        <string>-m</string>
        <string>http.server</string>
        <string>8765</string>
        <string>--bind</string>
        <string>127.0.0.1</string>
        <string>--directory</string>
        <string>$INSTALL_DIR</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/briefdesk-static.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/briefdesk-static.log</string>
</dict>
</plist>
EOF
chown "$INSTALL_USER" "$LAUNCHAGENTS_DIR/com.briefdesk.static.plist"

# Get Python version for user site-packages path
PYTHON_VERSION=$("$SYSTEM_PYTHON" -c "import sys; print(f'{sys.version_info.major}.{sys.version_info.minor}')" 2>/dev/null || echo "3.11")
USER_SITE_PACKAGES="$USER_HOME/.local/lib/python$PYTHON_VERSION/site-packages"

# Create LaunchAgent for search server
cat > "$LAUNCHAGENTS_DIR/com.briefdesk.server.plist" << EOF
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.briefdesk.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>$INSTALL_DIR/python3</string>
        <string>$INSTALL_DIR/search-server.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>$USER_SITE_PACKAGES</string>
        <key>HOME</key>
        <string>$USER_HOME</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/briefdesk-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/briefdesk-server.log</string>
</dict>
</plist>
EOF
chown "$INSTALL_USER" "$LAUNCHAGENTS_DIR/com.briefdesk.server.plist"

# Create LaunchAgent for search service (if node available)
if [ -n "$NODE_PATH" ]; then
    cat > "$LAUNCHAGENTS_DIR/com.briefdesk.search-service.plist" << EOF
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.briefdesk.search-service</string>
    <key>ProgramArguments</key>
    <array>
        <string>$DEVSAI_DIR/node</string>
        <string>$INSTALL_DIR/search-service.mjs</string>
    </array>
    <key>WorkingDirectory</key>
    <string>$INSTALL_DIR</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DEVSAI_LIB_PATH</key>
        <string>$DEVSAI_DIR/node_modules/devsai</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>$USER_HOME</string>
        <key>SEARCH_SERVICE_PORT</key>
        <string>19765</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/briefdesk-search-service.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/briefdesk-search-service.log</string>
</dict>
</plist>
EOF
    chown "$INSTALL_USER" "$LAUNCHAGENTS_DIR/com.briefdesk.search-service.plist"
fi

# Build gdrive-mcp if present
if [ -d "$INSTALL_DIR/gdrive-mcp" ] && [ -n "$NODE_PATH" ]; then
    echo "Building gdrive-mcp..."
    cd "$INSTALL_DIR/gdrive-mcp"
    sudo -u "$INSTALL_USER" npm install --silent 2>/dev/null || true
    sudo -u "$INSTALL_USER" npm run build --silent 2>/dev/null || true
fi

# Build gmail-mcp (custom read-only Gmail MCP)
if [ -d "$INSTALL_DIR/gmail-mcp" ] && [ -n "$NODE_PATH" ]; then
    echo "Building gmail-mcp (read-only)..."
    cd "$INSTALL_DIR/gmail-mcp"
    sudo -u "$INSTALL_USER" npm install --silent 2>/dev/null || true
    sudo -u "$INSTALL_USER" npm run build --silent 2>/dev/null || true
    
    # Verify build succeeded
    if [ -f "$INSTALL_DIR/gmail-mcp/dist/index.js" ]; then
        echo "  gmail-mcp built successfully"
    else
        echo "  WARNING: gmail-mcp build may have failed"
    fi
fi

# Remove any existing Gmail credentials with modify scope (security: re-auth with readonly)
if [ -f "$USER_HOME/.gmail-mcp/credentials.json" ]; then
    if grep -q "gmail.modify" "$USER_HOME/.gmail-mcp/credentials.json" 2>/dev/null; then
        rm -f "$USER_HOME/.gmail-mcp/credentials.json"
        echo "Removed Gmail credentials with modify scope - user will need to re-authenticate with read-only access"
    fi
fi

# Update .devsai.json to use the custom gmail-mcp (if config exists)
DEVSAI_CONFIG="$INSTALL_DIR/.devsai.json"
if [ -f "$DEVSAI_CONFIG" ] && [ -f "$INSTALL_DIR/gmail-mcp/dist/index.js" ]; then
    echo "Updating .devsai.json to use custom gmail-mcp..."
    # Use Python to safely update JSON (jq may not be available)
    sudo -u "$INSTALL_USER" "$SYSTEM_PYTHON" - "$DEVSAI_CONFIG" "$INSTALL_DIR" << 'PYEOF'
import json
import sys

config_path = sys.argv[1]
install_dir = sys.argv[2]

try:
    with open(config_path, 'r') as f:
        config = json.load(f)
    
    # Update gmail MCP to use the custom one
    if 'mcpServers' in config and 'gmail' in config['mcpServers']:
        config['mcpServers']['gmail'] = {
            "command": "node",
            "args": [f"{install_dir}/gmail-mcp/dist/index.js"]
        }
        
        with open(config_path, 'w') as f:
            json.dump(config, f, indent=2)
        print("  Updated gmail MCP config")
except Exception as e:
    print(f"  WARNING: Could not update .devsai.json: {e}")
PYEOF
fi

# Start services as the user
USER_ID=$(id -u "$INSTALL_USER")

# Unload any existing services first
sudo -u "$INSTALL_USER" launchctl bootout "gui/$USER_ID/com.briefdesk.static" 2>/dev/null || true
sudo -u "$INSTALL_USER" launchctl bootout "gui/$USER_ID/com.briefdesk.server" 2>/dev/null || true
sudo -u "$INSTALL_USER" launchctl bootout "gui/$USER_ID/com.briefdesk.search-service" 2>/dev/null || true

sleep 1

# Load services
sudo -u "$INSTALL_USER" launchctl bootstrap "gui/$USER_ID" "$LAUNCHAGENTS_DIR/com.briefdesk.static.plist" 2>/dev/null || true
sudo -u "$INSTALL_USER" launchctl bootstrap "gui/$USER_ID" "$LAUNCHAGENTS_DIR/com.briefdesk.server.plist" 2>/dev/null || true
if [ -f "$LAUNCHAGENTS_DIR/com.briefdesk.search-service.plist" ]; then
    sudo -u "$INSTALL_USER" launchctl bootstrap "gui/$USER_ID" "$LAUNCHAGENTS_DIR/com.briefdesk.search-service.plist" 2>/dev/null || true
fi

echo "BriefDesk installation complete!"
echo "Services are starting in the background..."

exit 0
