# pkg-installer/scripts/postinstall
#!/bin/bash
# BriefDesk Post-Install Script
# This runs as root after pkg installation, sets up all user-space components silently.
# Goal: User should have a FULLY working system after this script — no manual dependency installs.

set -e

# Get the user who invoked the installer (not root)
INSTALL_USER="${USER:-$(stat -f '%Su' /dev/console)}"
USER_HOME=$(eval echo ~$INSTALL_USER)

# Source files installed by pkg to /usr/local/share/briefdesk
PKG_SOURCE="/usr/local/share/briefdesk"

# User directories
INSTALL_DIR="$USER_HOME/.local/share/briefdesk"
DEVSAI_DIR="$USER_HOME/.local/share/devsai"
LAUNCHAGENTS_DIR="$USER_HOME/Library/LaunchAgents"

echo "Setting up BriefDesk for user: $INSTALL_USER (~2 minutes remaining)"

# Create directories
sudo -u "$INSTALL_USER" mkdir -p "$INSTALL_DIR"
sudo -u "$INSTALL_USER" mkdir -p "$DEVSAI_DIR"
sudo -u "$INSTALL_USER" mkdir -p "$LAUNCHAGENTS_DIR"
sudo -u "$INSTALL_USER" mkdir -p "$USER_HOME/.local/bin"

# Copy application files
cp -R "$PKG_SOURCE/"* "$INSTALL_DIR/"
chown -R "$INSTALL_USER" "$INSTALL_DIR"

# ============================================================
# 1. PYTHON — detect or trigger install
# ============================================================
echo "Step 1/6: Detecting Python... (~2 minutes remaining)"
SYSTEM_PYTHON=""

check_python() {
    local py="$1"
    if [ -x "$py" ] && "$py" -c "import sys; exit(0 if sys.version_info >= (3, 9) else 1)" 2>/dev/null; then
        return 0
    fi
    return 1
}

# Check all common Python locations (prefer Homebrew, then Framework, then system)
for py in \
    /opt/homebrew/bin/python3 \
    /opt/homebrew/bin/python3.{14,13,12,11,10,9} \
    /Library/Frameworks/Python.framework/Versions/{3.14,3.13,3.12,3.11,3.10,3.9}/bin/python3 \
    /usr/local/bin/python3 \
    /usr/local/bin/python3.{14,13,12,11,10,9} \
    "$USER_HOME/.pyenv/shims/python3" \
    /usr/bin/python3; do
    if [ -z "$SYSTEM_PYTHON" ] && [ -f "$py" ] && check_python "$py"; then
        SYSTEM_PYTHON="$py"
        break
    fi
done

# Last resort: which
if [ -z "$SYSTEM_PYTHON" ]; then
    FOUND_PYTHON=$(which python3 2>/dev/null || true)
    if [ -n "$FOUND_PYTHON" ] && check_python "$FOUND_PYTHON"; then
        SYSTEM_PYTHON="$FOUND_PYTHON"
    fi
fi

# If still not found, trigger Xcode CLI Tools install (includes Python 3)
if [ -z "$SYSTEM_PYTHON" ]; then
    echo "Python 3.9+ not found. Triggering Xcode Command Line Tools install..."
    # This shows a system dialog asking the user to install
    sudo -u "$INSTALL_USER" xcode-select --install 2>/dev/null || true
    # Wait up to 5 minutes for Python to appear (user needs to click Install)
    for i in $(seq 1 60); do
        sleep 5
        if [ -f "/usr/bin/python3" ] && check_python "/usr/bin/python3"; then
            SYSTEM_PYTHON="/usr/bin/python3"
            echo "Python installed via Xcode CLI Tools"
            break
        fi
    done
fi

# Final fallback
if [ -z "$SYSTEM_PYTHON" ]; then
    echo "WARNING: Could not find Python 3.9+. History search may not work."
    SYSTEM_PYTHON="/usr/bin/python3"
fi

# Resolve symlinks to get the real binary path (needed for FDA)
PYTHON_REAL_PATH=$("$SYSTEM_PYTHON" -c "import os, sys; print(os.path.realpath(sys.executable))" 2>/dev/null || echo "$SYSTEM_PYTHON")
PYTHON_VERSION=$("$SYSTEM_PYTHON" --version 2>&1 || echo "unknown")
echo "Using Python: $SYSTEM_PYTHON → $PYTHON_REAL_PATH ($PYTHON_VERSION)"

# Create Python wrapper
cat > "$INSTALL_DIR/python3" << PYEOF
#!/bin/bash
exec "$SYSTEM_PYTHON" "\$@"
PYEOF
chmod +x "$INSTALL_DIR/python3"
chown "$INSTALL_USER" "$INSTALL_DIR/python3"

# Install Python dependencies (silently)
echo "Step 1/6: Installing Python dependencies... (~2 minutes remaining)"
if ! sudo -u "$INSTALL_USER" "$SYSTEM_PYTHON" -m pip --version >/dev/null 2>&1; then
    sudo -u "$INSTALL_USER" "$SYSTEM_PYTHON" -m ensurepip --user 2>/dev/null || true
fi
sudo -u "$INSTALL_USER" "$SYSTEM_PYTHON" -m pip install --user --quiet google-auth-oauthlib google-api-python-client 2>/dev/null || \
    sudo -u "$INSTALL_USER" "$SYSTEM_PYTHON" -m pip install --break-system-packages --user --quiet google-auth-oauthlib google-api-python-client 2>/dev/null || \
    echo "WARNING: Could not install Google API libraries."

# ============================================================
# 2. NODE.JS — detect or auto-download
# ============================================================
echo "Step 2/6: Setting up Node.js... (~90 seconds remaining)"
NODE_PATH=""
NODE_SOURCE="system"

# Check system Node installations
for np in /opt/homebrew/bin/node /usr/local/bin/node /opt/homebrew/opt/node/bin/node; do
    if [ -z "$NODE_PATH" ] && [ -f "$np" ]; then
        NODE_PATH="$np"
        break
    fi
done

# If not found, auto-download Node.js LTS
if [ -z "$NODE_PATH" ]; then
    echo "Node.js not found. Downloading Node.js LTS..."
    NODE_SOURCE="downloaded"
    ARCH=$(uname -m)
    if [ "$ARCH" = "arm64" ]; then
        NODE_ARCH="arm64"
    else
        NODE_ARCH="x64"
    fi
    NODE_DL_URL="https://nodejs.org/dist/v22.13.1/node-v22.13.1-darwin-${NODE_ARCH}.tar.gz"
    NODE_TMP="/tmp/briefdesk-node-install"
    mkdir -p "$NODE_TMP"
    if curl -fsSL "$NODE_DL_URL" -o "$NODE_TMP/node.tar.gz" 2>/dev/null; then
        tar -xzf "$NODE_TMP/node.tar.gz" -C "$NODE_TMP" --strip-components=1 2>/dev/null
        if [ -f "$NODE_TMP/bin/node" ]; then
            NODE_PATH="$NODE_TMP/bin/node"
            # Copy full Node installation to persistent location
            # so npm works for MCP builds and beyond the install session
            cp "$NODE_TMP/bin/node" "$DEVSAI_DIR/node" 2>/dev/null || true
            if [ -d "$NODE_TMP/lib/node_modules/npm" ]; then
                # Copy npm's lib so npm-cli.js can resolve its require() paths
                cp -r "$NODE_TMP/lib" "$DEVSAI_DIR/" 2>/dev/null || true
                # Create npm and npx wrappers that use our node binary
                cat > "$DEVSAI_DIR/npm" << NPMEOF
#!/bin/bash
exec "$DEVSAI_DIR/node" "$DEVSAI_DIR/lib/node_modules/npm/bin/npm-cli.js" "\$@"
NPMEOF
                cat > "$DEVSAI_DIR/npx" << NPXEOF
#!/bin/bash
exec "$DEVSAI_DIR/node" "$DEVSAI_DIR/lib/node_modules/npm/bin/npx-cli.js" "\$@"
NPXEOF
            fi
            chmod +x "$DEVSAI_DIR/node" "$DEVSAI_DIR/npm" "$DEVSAI_DIR/npx" 2>/dev/null || true
            chown -R "$INSTALL_USER" "$DEVSAI_DIR" 2>/dev/null || true
            # Point NODE_PATH to the persistent copy (temp dir is cleaned up below)
            NODE_PATH="$DEVSAI_DIR/node"
            echo "  Node.js downloaded successfully"
        fi
    else
        echo "WARNING: Failed to download Node.js (no internet?)"
    fi
    rm -rf "$NODE_TMP"
fi

# Find npm
NPM_PATH=""
if [ -f "/usr/local/bin/npm" ]; then
    NPM_PATH="/usr/local/bin/npm"
elif [ -f "/opt/homebrew/bin/npm" ]; then
    NPM_PATH="/opt/homebrew/bin/npm"
elif [ -n "$NODE_PATH" ] && [ -f "$(dirname "$NODE_PATH")/npm" ]; then
    NPM_PATH="$(dirname "$NODE_PATH")/npm"
elif [ -f "$DEVSAI_DIR/npm" ]; then
    # Use the npm we downloaded/copied alongside node
    NPM_PATH="$DEVSAI_DIR/npm"
else
    NPM_PATH="$(command -v npm 2>/dev/null || true)"
fi

if [ -n "$NPM_PATH" ]; then
    echo "Using npm: $NPM_PATH"
else
    echo "WARNING: npm not found — MCP server builds will be skipped."
fi

# Copy Node binary to persistent location (skip if already copied during download)
if [ -n "$NODE_PATH" ] && [ ! -f "$DEVSAI_DIR/node" ]; then
    cp "$NODE_PATH" "$DEVSAI_DIR/node"
    chmod +x "$DEVSAI_DIR/node"
    chown "$INSTALL_USER" "$DEVSAI_DIR/node"
fi
if [ -f "$DEVSAI_DIR/node" ]; then
    NODE_VERSION=$("$DEVSAI_DIR/node" --version 2>/dev/null || echo "unknown")
    echo "Using Node: $NODE_PATH ($NODE_VERSION, source: $NODE_SOURCE)"
fi

# Build a PATH that includes the node binary's directory.
# npm's shebang is #!/usr/bin/env node — when running via sudo, the PATH is sanitized
# and may not include /opt/homebrew/bin, causing "env: node: No such file or directory".
NODE_BIN_DIR=""
if [ -n "$NODE_PATH" ]; then
    NODE_BIN_DIR="$(dirname "$NODE_PATH")"
fi
SUDO_PATH="${NODE_BIN_DIR:+$NODE_BIN_DIR:}/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin"

# ============================================================
# 3. DEVSAI CLI — install from bundled tarball
# ============================================================
echo "Step 3/6: Installing search engine... (~75 seconds remaining)"
DEVSAI_TARBALL="$INSTALL_DIR/devsai/devsai.tgz"
DEVSAI_INSTALLED=false

if [ -f "$DEVSAI_TARBALL" ] && [ -n "$NPM_PATH" ]; then
    echo "Installing devsai CLI from tarball..."
    # HOME must be set for npm cache/config; PATH must include node binary dir for npm's shebang
    sudo -u "$INSTALL_USER" HOME="$USER_HOME" PATH="$SUDO_PATH" "$NPM_PATH" install --prefix "$DEVSAI_DIR" "$DEVSAI_TARBALL" 2>&1 | tail -5 || \
        echo "WARNING: devsai tarball install failed."
    if [ -f "$DEVSAI_DIR/node_modules/devsai/dist/index.js" ]; then
        DEVSAI_INSTALLED=true
        echo "   ✓ Installed devsai from tarball"
    else
        echo "   ⚠ devsai tarball install did not produce expected files"
    fi
else
    [ ! -f "$DEVSAI_TARBALL" ] && echo "   ⚠ devsai tarball not found at $DEVSAI_TARBALL"
    [ -z "$NPM_PATH" ] && echo "   ⚠ npm not available — cannot install devsai"
fi

# Create wrapper script
if [ "$DEVSAI_INSTALLED" = true ]; then
    cat > "$DEVSAI_DIR/devsai.sh" << 'EOF'
#!/bin/bash
exec "$HOME/.local/share/devsai/node" "$HOME/.local/share/devsai/node_modules/devsai/dist/index.js" "$@"
EOF
    chmod +x "$DEVSAI_DIR/devsai.sh"
    chown "$INSTALL_USER" "$DEVSAI_DIR/devsai.sh"
    sudo -u "$INSTALL_USER" ln -sf "$DEVSAI_DIR/devsai.sh" "$USER_HOME/.local/bin/devsai"
fi

# ============================================================
# 4. LAUNCH AGENTS
# ============================================================
echo "Step 4/6: Configuring services... (~60 seconds remaining)"
PYTHON_VER=$("$SYSTEM_PYTHON" -c "import sys; print(f'{sys.version_info.major}.{sys.version_info.minor}')" 2>/dev/null || echo "3.11")
USER_SITE_PACKAGES="$USER_HOME/.local/lib/python$PYTHON_VER/site-packages"

# Static file server
cat > "$LAUNCHAGENTS_DIR/com.briefdesk.static.plist" << EOF
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.briefdesk.static</string>
    <key>ProgramArguments</key>
    <array>
        <string>$INSTALL_DIR/python3</string>
        <string>-m</string>
        <string>http.server</string>
        <string>8765</string>
        <string>--bind</string>
        <string>127.0.0.1</string>
        <string>--directory</string>
        <string>$INSTALL_DIR</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/briefdesk-static.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/briefdesk-static.log</string>
</dict>
</plist>
EOF
chown "$INSTALL_USER" "$LAUNCHAGENTS_DIR/com.briefdesk.static.plist"

# API server
cat > "$LAUNCHAGENTS_DIR/com.briefdesk.server.plist" << EOF
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.briefdesk.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>$INSTALL_DIR/python3</string>
        <string>$INSTALL_DIR/search-server.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>$USER_SITE_PACKAGES</string>
        <key>HOME</key>
        <string>$USER_HOME</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/briefdesk-server.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/briefdesk-server.log</string>
</dict>
</plist>
EOF
chown "$INSTALL_USER" "$LAUNCHAGENTS_DIR/com.briefdesk.server.plist"

# Search service (Node.js)
if [ -n "$NODE_PATH" ]; then
    cat > "$LAUNCHAGENTS_DIR/com.briefdesk.search-service.plist" << EOF
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.briefdesk.search-service</string>
    <key>ProgramArguments</key>
    <array>
        <string>$DEVSAI_DIR/node</string>
        <string>$INSTALL_DIR/search-service.mjs</string>
    </array>
    <key>WorkingDirectory</key>
    <string>$INSTALL_DIR</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DEVSAI_LIB_PATH</key>
        <string>$DEVSAI_DIR/node_modules/devsai</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>$USER_HOME</string>
        <key>SEARCH_SERVICE_PORT</key>
        <string>19765</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/briefdesk-search-service.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/briefdesk-search-service.log</string>
</dict>
</plist>
EOF
    chown "$INSTALL_USER" "$LAUNCHAGENTS_DIR/com.briefdesk.search-service.plist"
fi

# ============================================================
# 5. GITHUB MCP SERVER — download binary
# ============================================================
echo "Step 5/6: Downloading integrations... (~45 seconds remaining)"
GITHUB_MCP_BIN="$INSTALL_DIR/github-mcp-server"
if [ ! -f "$GITHUB_MCP_BIN" ]; then
    ARCH=$(uname -m)
    GH_MCP_ASSET="github-mcp-server_Darwin_${ARCH}.tar.gz"
    [ "$ARCH" != "arm64" ] && GH_MCP_ASSET="github-mcp-server_Darwin_x86_64.tar.gz"
    if curl -fsSL "https://github.com/github/github-mcp-server/releases/latest/download/$GH_MCP_ASSET" -o "/tmp/$GH_MCP_ASSET" 2>/dev/null; then
        tar -xzf "/tmp/$GH_MCP_ASSET" -C /tmp github-mcp-server 2>/dev/null
        if [ -f "/tmp/github-mcp-server" ]; then
            mv "/tmp/github-mcp-server" "$GITHUB_MCP_BIN"
            chmod +x "$GITHUB_MCP_BIN"
            chown "$INSTALL_USER" "$GITHUB_MCP_BIN"
        fi
        rm -f "/tmp/$GH_MCP_ASSET"
    fi
fi

# ============================================================
# 6. BUILD MCP SERVERS (gdrive, gmail)
# ============================================================
echo "Step 6/6: Building integrations... (~30 seconds remaining)"
for MCP_NAME in gdrive-mcp gmail-mcp; do
    if [ -d "$INSTALL_DIR/$MCP_NAME" ] && [ -n "$NPM_PATH" ]; then
        cd "$INSTALL_DIR/$MCP_NAME"
        # HOME for npm cache/config; PATH for npm's #!/usr/bin/env node shebang; --include=dev for tsc
        sudo -u "$INSTALL_USER" HOME="$USER_HOME" PATH="$SUDO_PATH" "$NPM_PATH" install --include=dev 2>&1 | tail -3 || echo "WARNING: $MCP_NAME npm install failed"
        sudo -u "$INSTALL_USER" HOME="$USER_HOME" PATH="$SUDO_PATH" "$NPM_PATH" run build 2>&1 | tail -3 || echo "WARNING: $MCP_NAME build failed"
        if [ -f "$INSTALL_DIR/$MCP_NAME/dist/index.js" ]; then
            echo "   ✓ $MCP_NAME built successfully"
        else
            echo "   ⚠ $MCP_NAME build failed - will retry on next install"
        fi
    fi
done

# Remove old Gmail credentials with modify scope (security)
if [ -f "$USER_HOME/.gmail-mcp/credentials.json" ]; then
    grep -q "gmail.modify" "$USER_HOME/.gmail-mcp/credentials.json" 2>/dev/null && \
        rm -f "$USER_HOME/.gmail-mcp/credentials.json"
fi

# Update .devsai.json gmail config if needed
DEVSAI_CONFIG="$INSTALL_DIR/.devsai.json"
if [ -f "$DEVSAI_CONFIG" ] && [ -f "$INSTALL_DIR/gmail-mcp/dist/index.js" ]; then
    sudo -u "$INSTALL_USER" "$SYSTEM_PYTHON" - "$DEVSAI_CONFIG" "$INSTALL_DIR" << 'PYEOF'
import json, sys
config_path, install_dir = sys.argv[1], sys.argv[2]
try:
    with open(config_path, 'r') as f: config = json.load(f)
    if 'mcpServers' in config and 'gmail' in config['mcpServers']:
        config['mcpServers']['gmail'] = {"command": "node", "args": [f"{install_dir}/gmail-mcp/dist/index.js"]}
        with open(config_path, 'w') as f: json.dump(config, f, indent=2)
except: pass
PYEOF
fi

# ============================================================
# 7. FIREFOX EXTENSION AUTO-INSTALL
# ============================================================
FIREFOX_EXT_SRC="$INSTALL_DIR/browser-extension/firefox"
if [ -d "$FIREFOX_EXT_SRC" ]; then
    # Find Firefox profile directory
    FIREFOX_PROFILES="$USER_HOME/Library/Application Support/Firefox/Profiles"
    if [ -d "$FIREFOX_PROFILES" ]; then
        for profile_dir in "$FIREFOX_PROFILES"/*.default-release "$FIREFOX_PROFILES"/*.default; do
            if [ -d "$profile_dir" ]; then
                EXT_DIR="$profile_dir/extensions"
                sudo -u "$INSTALL_USER" mkdir -p "$EXT_DIR"
                # Copy as XPI (zip) or use pointer file
                # Firefox reads a file named with the extension ID containing the path
                echo "$FIREFOX_EXT_SRC" > "$EXT_DIR/newtab@briefdesk.local"
                chown "$INSTALL_USER" "$EXT_DIR/newtab@briefdesk.local"
                echo "  Firefox extension linked in $(basename "$profile_dir")"
            fi
        done
    fi
fi

# ============================================================
# 8. SAVE SYSTEM METADATA (for installer.html to read)
# ============================================================
cat > "$INSTALL_DIR/.install-meta.json" << METAEOF
{
  "python_path": "$SYSTEM_PYTHON",
  "python_real_path": "$PYTHON_REAL_PATH",
  "python_version": "$PYTHON_VERSION",
  "node_path": "${NODE_PATH:-}",
  "node_version": "${NODE_VERSION:-}",
  "node_source": "$NODE_SOURCE",
  "node_fda_path": "$DEVSAI_DIR/node",
  "install_dir": "$INSTALL_DIR",
  "devsai_dir": "$DEVSAI_DIR",
  "github_mcp": "$([ -f "$GITHUB_MCP_BIN" ] && echo 'installed' || echo 'missing')",
  "devsai_cli": "$([ -f "$DEVSAI_DIR/devsai.sh" ] && echo 'installed' || echo 'missing')",
  "installed_at": "$(date -u +%Y-%m-%dT%H:%M:%SZ)"
}
METAEOF
chown "$INSTALL_USER" "$INSTALL_DIR/.install-meta.json"

# ============================================================
# 9. START SERVICES
# ============================================================
echo "Starting services... (almost done)"
USER_ID=$(id -u "$INSTALL_USER")

# Unload existing
sudo -u "$INSTALL_USER" launchctl bootout "gui/$USER_ID/com.briefdesk.static" 2>/dev/null || true
sudo -u "$INSTALL_USER" launchctl bootout "gui/$USER_ID/com.briefdesk.server" 2>/dev/null || true
sudo -u "$INSTALL_USER" launchctl bootout "gui/$USER_ID/com.briefdesk.search-service" 2>/dev/null || true
sleep 1

# Load new
sudo -u "$INSTALL_USER" launchctl bootstrap "gui/$USER_ID" "$LAUNCHAGENTS_DIR/com.briefdesk.static.plist" 2>/dev/null || true
sudo -u "$INSTALL_USER" launchctl bootstrap "gui/$USER_ID" "$LAUNCHAGENTS_DIR/com.briefdesk.server.plist" 2>/dev/null || true
if [ -f "$LAUNCHAGENTS_DIR/com.briefdesk.search-service.plist" ]; then
    sudo -u "$INSTALL_USER" launchctl bootstrap "gui/$USER_ID" "$LAUNCHAGENTS_DIR/com.briefdesk.search-service.plist" 2>/dev/null || true
fi

echo "BriefDesk installation complete!"
exit 0
