# launchagents/com.startpage.search.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.startpage.search</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>HOME_DIR/.local/share/safari_start_page/search-server.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SLACK_WORKSPACE</key>
        <string>your-workspace</string>
        <key>ATLASSIAN_DOMAIN</key>
        <string>your-domain.atlassian.net</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/startpage-search.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/startpage-search.log</string>
</dict>
</plist>
